# manage/postgres/postgres_backup.timer
[Unit]
Description=postgres backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
